# v2/systemd/v2-stack.service
[Unit]
Description=Auto Trader V2 web-first stack (control API + operator web)
After=network-online.target
Wants=network-online.target

[Service]
Type=simple
User=bot
WorkingDirectory=/home/bot/autotrade/auto_trade_V2
Environment=PYTHONUNBUFFERED=1
EnvironmentFile=/home/bot/autotrade/auto_trade_V2/.env
ExecStart=/usr/bin/env bash /home/bot/autotrade/auto_trade_V2/v2/scripts/run_stack.sh --profile ra_2026_alpha_v2_expansion_verified_q070 --mode live --env prod --env-file .env --host 127.0.0.1 --port 8101 --operator-web --no-discord-bot
Restart=always
RestartSec=2
KillMode=control-group
TimeoutStopSec=10

[Install]
WantedBy=multi-user.target
